# build-tools/debian/postinst
#!/bin/bash

chmod +x /etc/init.d/igs
update-rc.d igs defaults
/etc/init.d/igs start
